# packaging/monasca-agent-deb/postrm
#! /bin/bash
set -e

if [ "$1" = purge ]; then
    update-rc.d monasca-agent remove
    deluser mon-agent
fi

exit 0
