# deb/init/mon-agent.conf
# Startup script for the Mon-Agent

description "Mon-Agent application"
start on runlevel [2345]

console log
respawn

setgid mon_agent
setuid mon_agent
exec /usr/bin/python /opt/mon/agent.py start
