# monagent/packaging/datadog-agent-deb/postrm
#! /bin/bash
set -e

if [ "$1" = purge ]; then
    update-rc.d mon-agent remove
    deluser mon-agent
fi

exit 0
